# profile.d/brew-overrides.sh.macos

use_brew_opt_override() {
	local name="$1"

	homedir_path_prepend "/usr/local/opt/$name/libexec/bin" ||\
	homedir_path_prepend "/usr/local/opt/$name/bin" ||\
	homedir_path_prepend "/usr/local/opt/$name/libexec"

	homedir_path_prepend  PKG_CONFIG_PATH "/usr/local/opt/$name/lib/pkgconfig"
	return 0
}

use_brew_opt_override python
use_brew_opt_override ccache

# use_brew_opt_override curl

