# profile.d/brew-overrides.sh.macos
#!/usr/bin/env bash

use_brew_opt_override() {
	local name="$1"

	shpath_pre "/usr/local/opt/$name/libexec/bin" || \
	shpath_pre "/usr/local/opt/$name/bin" || \
	shpath_pre "/usr/local/opt/$name/libexec"

	shpath_pre  PKG_CONFIG_PATH "/usr/local/opt/$name/lib/pkgconfig"
	return 0
}

use_brew_opt_override python
use_brew_opt_override ccache

# use_brew_opt_override curl

